# crontab-config/cvndb-pull.cron
*/5 * * * *	root	/srv/cvn/git/infrastructure/bin/cvndb-pull > /dev/null
