# crontab-config/cvndb-CVNBot14-publish.cron
*/5 * * * *	root	/srv/cvn/git/infrastructure/bin/cvndb-CVNBot14-publish > /dev/null
