# crontab-config/cvndb-CVNBot14-publish.cron
*/5 * * * *	root	/data/project/cvn-common/git/infrastructure/bin/cvndb-CVNBot14-publish > /dev/null
