# crontab-config/cvndb-pull.cron
*/5 * * * *	root	/data/project/cvn-common/git/infrastructure/bin/cvndb-pull > /dev/null
