# crontab-config/stillalive.cron
*/5 * * * *     root    /data/project/cvn-common/git/stillalive/bin/stillalive --pool=$(hostname) > /dev/null
