# crontab-config/stillalive.cron
*/2 * * * *	root	/srv/cvn/git/stillalive/bin/stillalive --pool=$(hostname) > /dev/null
